# services/scheduler.service
[Unit]
Description=Pisces Scheduler Service
After=network.target hwcontrol.service
#StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
#User=pi
WorkingDirectory=/home/pi/Repositories/pisces/
ExecStart=/usr/bin/env python3 /home/pi/Repositories/pisces/scheduler/scheduler.py
Environment="PYTHONPATH=/home/pi/Repositories/pisces/shared"

#StandardOutput=
#StandatdError=inherit

[Install]
WantedBy=multi-user.target